# deploy/tryallfun-panel.service
[Unit]
Description=TryAllFun Server Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/wwwroot/Kunpanel.456.life
ExecStart=/home/wwwroot/Kunpanel.456.life/tryallfun-panel
Environment=TAF_ADDR=127.0.0.1:8088
Environment=TAF_DATA_DIR=/var/lib/tryallfun-panel
Environment=TAF_FILE_ROOT=/home/wwwroot
Environment=TAF_PROJECT_DIR=/home/wwwroot/Kunpanel.456.life
Environment=TAF_BINARY_PATH=/home/wwwroot/Kunpanel.456.life/tryallfun-panel
Environment=TAF_NGINX_BIN=/usr/local/nginx/sbin/nginx
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=false
ProtectHome=false
ProtectSystem=false

[Install]
WantedBy=multi-user.target
